# deploy/zeus-hermes.timer
[Unit]
Description=Zeus Hermes — briefing alle 08:05 UTC (dopo orchestrator)
Requires=zeus-hermes.service

[Timer]
# 08:05 UTC — 35 minuti dopo l'orchestrator, report gia' pronto
OnCalendar=*-*-* 08:05:00 UTC
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
