# deploy/zeus-hermes.service
[Unit]
Description=Zeus Hermes — briefing mattutino via Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=zeus
WorkingDirectory=/opt/zeus
EnvironmentFile=/opt/zeus/.env
ExecStart=/opt/zeus/venv/bin/python hermes/briefing.py
StandardOutput=append:/opt/zeus/log/hermes.log
StandardError=append:/opt/zeus/log/hermes.log

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
