# hermes/hermes-bot.service
[Unit]
Description=Hermes Zeus Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zeus/hermes
EnvironmentFile=/opt/zeus/hermes/.env
ExecStart=/opt/zeus/.venv/bin/python3 /opt/zeus/hermes/bot.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/hermes_bot.log
StandardError=append:/var/log/hermes_bot.log

[Install]
WantedBy=multi-user.target
